# deploy/systemd/hermes-team.target
[Unit]
Description=Hermes Team Agent Gateway Target
After=network.target

[Install]
WantedBy=multi-user.target
